# crontab
# crypto-monitor supercronic crontab
# Formato: minuto hora dia mês diasemana comando

# Hyperliquid: funding rates + posições (a cada 15 min)
*/15 * * * * python /app/agents/hyperliquid.py

# MEXC: carteira + BTC price (a cada 15 min)
*/15 * * * * python /app/agents/mexc.py

# News + Sentimento (a cada 10 min)
*/10 * * * * python /app/agents/sentiment.py

# Polymarket scanner (a cada 30 min)
*/30 * * * * python /app/agents/polymarket.py

# Consensus orchestrator (a cada 4h, na hora cheia)
0 */4 * * * python /app/agents/orchestrator.py

# Kronos: previsão + gráficos 1H/4H/Diário → Telegram [KRONOS] (a cada 4h, +15min do orchestrator)
15 */4 * * * python /app/vps/kronos_signal.py >> /data/kronos.log 2>&1

# Kronos scorecard: avalia previsões maduras (Telegram só se houver novidades)
30 */4 * * * python /app/vps/kronos_scorecard.py >> /data/kronos_scorecard.log 2>&1
# Relatório completo semanal (domingo 12:30 UTC)
30 12 * * 0 python /app/vps/kronos_scorecard.py --weekly >> /data/kronos_scorecard.log 2>&1
