# crontab
# crypto-monitor (Railway) — sem Kronos
# Kronos: Hetzner BTCCURSOR — vps/BTCCURSOR.md

# Hyperliquid: funding rates + posições (a cada 15 min)
*/15 * * * * python /app/agents/hyperliquid.py

# MEXC: carteira + BTC price (a cada 15 min)
*/15 * * * * python /app/agents/mexc.py

# News + Sentimento (a cada 10 min)
*/10 * * * * python /app/agents/sentiment.py

# Polymarket scanner (a cada 30 min)
*/30 * * * * python /app/agents/polymarket.py

# Consensus orchestrator (a cada 4h, na hora cheia)
0 */4 * * * python /app/agents/orchestrator.py
