# arquivo-kronos/codigo/crontab
# crypto-monitor supercronic crontab
# Formato: minuto hora dia mês diasemana comando

# Hyperliquid: funding rates + posições (a cada 15 min)
*/15 * * * * python /app/agents/hyperliquid.py

# MEXC: carteira + BTC price (a cada 15 min)
*/15 * * * * python /app/agents/mexc.py

# News + Sentimento (a cada 10 min)
*/10 * * * * python /app/agents/sentiment.py

# Polymarket scanner (a cada 30 min)
*/30 * * * * python /app/agents/polymarket.py

# Consensus orchestrator (a cada 4h, na hora cheia)
0 */4 * * * python /app/agents/orchestrator.py

# Kronos: alerta no fechamento de cada candle MEXC/UTC (lock anti-OOM)
2 * * * * /app/vps/kronos_candle_cron.sh /data/kronos.log
5 0 * * * /app/vps/kronos_run.sh /data/kronos.log 1d
# Watchdog: avisa se ficar >90min sem alerta OK (runs 1H horários)
50 * * * * python /app/vps/kronos_watchdog.py >> /data/kronos_watchdog.log 2>&1

# Kronos scorecard: após runs 4H (tempo para o modelo terminar)
45 0,4,8,12,16,20 * * * python /app/vps/kronos_scorecard.py >> /data/kronos_scorecard.log 2>&1
# Relatório completo semanal (domingo 12:30 UTC)
30 12 * * 0 python /app/vps/kronos_scorecard.py --weekly >> /data/kronos_scorecard.log 2>&1

# Kronos: relatório diário no Telegram — 12:00 horário de Brasília (15:00 UTC)
0 15 * * * python /app/vps/kronos_daily_report.py >> /data/kronos_daily_report.log 2>&1
